# resources/amazon/deploy/charms/trusty/deploy-node/templates/upstart.conf
description "deploy-node"
author "root <root@ubuntu-trusty-1.c.jujutest-1158.internal>"

start on runlevel [2345]
stop on runlevel [016]

respawn

console log
script
    echo Fake service; sleeping for an hour...
    sleep 360
end script
